# bootstrap/service/configs/dev/supervisord.conf
[program:code-changes-watcher]

command=%(ENV_BOOTSTRAP_DIR)s/service/scripts/dev/code-changes-watcher.sh ; command to run
priority=2
autostart=true
redirect_stderr=true
autorestart=true
