# bootstrap/service/service-process/supervisord.conf
[unix_http_server]
chown=%(ENV_APP_USER)s
file=/tmp/supervisor.sock
username=%(ENV_APP_USER)s
password=%(ENV_APP_USER)s

[supervisord]
user=%(ENV_APP_USER)s
logfile=/dev/null
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
username=%(ENV_APP_USER)s
password=%(ENV_APP_USER)s

[program:process-tasks-general]
command=celery --app=green_detective worker --beat --autoscale=%(ENV_CELERYD_MAX_CONCURRENCY_GENERAL)s,%(ENV_CELERYD_MIN_CONCURRENCY_GENERAL)s --hostname=gd_general --loglevel=info --queues=gd_general -E --scheduler django_celery_beat.schedulers:DatabaseScheduler
priority=100
autostart=true
autorestart=false
redirect_stderr=true
stopsignal=KILL
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
directory=%(ENV_WORKDIR)s/ ; command runing dir

[program:process-tasks-scrape]
command=celery --app=green_detective worker --autoscale=%(ENV_CELERYD_MAX_CONCURRENCY_SCRAPE)s,%(ENV_CELERYD_MIN_CONCURRENCY_SCRAPE)s --hostname=gd_scrape --loglevel=info --queues=gd_scrape -E
priority=100
autostart=true
autorestart=false
redirect_stderr=true
stopsignal=KILL
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
directory=%(ENV_WORKDIR)s/ ; command runing dir

[program:process-tasks-pre-staging]
command=celery --app=green_detective worker --autoscale=%(ENV_CELERYD_MAX_CONCURRENCY_PRE_STAGING)s,%(ENV_CELERYD_MIN_CONCURRENCY_PRE_STAGING)s --hostname=gd_pre_staging --loglevel=info --queues=gd_pre_staging -E
priority=100
autostart=true
autorestart=false
redirect_stderr=true
stopsignal=KILL
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
directory=%(ENV_WORKDIR)s/ ; command runing dir

[program:process-tasks-post-staging]
command=celery --app=green_detective worker --autoscale=%(ENV_CELERYD_MAX_CONCURRENCY_POST_STAGING)s,%(ENV_CELERYD_MIN_CONCURRENCY_POST_STAGING)s --hostname=gd_post_staging --loglevel=info --queues=gd_post_staging -E
priority=100
autostart=true
autorestart=false
redirect_stderr=true
stopsignal=KILL
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
directory=%(ENV_WORKDIR)s/ ; command runing dir
